# host/com.claude.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>__HOST_DIR__/claude_dashboard_host.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOST_DIR__/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>__HOST_DIR__/dashboard.log</string>
</dict>
</plist>
